# infra/ansible/roles/recorder/templates/adsb-recorder.service.j2
# ==============================================================================
# File: infra/ansible/roles/recorder/templates/adsb-recorder.service.j2
# Service: adsb-recorder.service
# Version: 0.3.9
# Managed By: Ansible
# ==============================================================================

[Unit]
Description=ADSB Raw Data Recorder
After=docker.service network.target
Requires=docker.service

[Service]
Type=simple
# Updated path for v0.3.9 (System Binary Path)
ExecStart=/usr/local/bin/adsb_record.sh
Restart=always
RestartSec=10
User=root
Group=root

[Install]
WantedBy=multi-user.target
