# infra/ansible/roles/recorder/templates/gnss-recorder.service.j2
# ==============================================================================
# File: infra/ansible/roles/recorder/templates/gnss-recorder.service.j2
# Service: gnss-recorder.service
# Version: 0.3.9
# Managed By: Ansible
# ==============================================================================

[Unit]
Description=GNSS Data Recorder Service
After=network.target
# Conditional logic: Wait for the correct GNSS source (RTK vs Standard)
{% if ntrip_host is defined %}
Requires=gnss-receiver.service
After=gnss-receiver.service
{% else %}
Requires=gpsd.service
After=gpsd.service
{% endif %}

[Service]
Type=simple
# Updated path for v0.3.9 (System Binary Path)
ExecStart=/usr/local/bin/gnss_record.sh
Restart=always
RestartSec=10
User=root
Group=root

[Install]
WantedBy=multi-user.target
