# infra/ansible/roles/recorder/templates/smart-logger.service.j2
# ==============================================================================
# FILE: /etc/systemd/system/smart-logger.service
# VERSION: 0.4.5 (Managed by Ansible)
# ==============================================================================
# DESCRIPTION:
#   Systemd unit for the ADSB Smart Logger.
#   Ensures the service runs as 'pi', starts after the network/GPSD,
#   and injects environment variables for Hardware Agnosticism.
# ==============================================================================

[Unit]
Description=Smart ADS-B Data Logger (Academic v1.4.0)
After=network.target readsb.service gpsd.service
Wants=gpsd.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/adsb-sensor

# -----------------------------------------------------------
# HARDWARE CONFIGURATION
# -----------------------------------------------------------
# These variables allow the same Python script to run on
# Sensor-North (HAT) and Sensor-West (USB) without code changes.
Environment="SENSOR_ID={{ sensor_id }}"
Environment="SENSOR_GNSS_PATH={{ gnss_device_path | default('/dev/ttyUSB0') }}"

# -----------------------------------------------------------
# RUNTIME CONFIGURATION
# -----------------------------------------------------------
# Force Python to unbuffer output so logs appear instantly in journalctl
Environment="PYTHONUNBUFFERED=1"

# EXECUTION PATH (Standardized to /opt)
ExecStart=/usr/bin/python3 /opt/adsb-sensor/smart_adsb_logger.py

# RESTART POLICY
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
