# infra/ansible/roles/recorder/templates/smart-logger.service.j2
# ==============================================================================
# FILE: /etc/systemd/system/smart-logger.service
# VERSION: 1.0.0 (Research Edition)
# ==============================================================================
# DESCRIPTION:
#   Systemd unit for the ADSB Smart Logger.
#   Updated to support dynamic aircraft logging for research.
# ==============================================================================

[Unit]
Description=Smart ADS-B Data Logger (Research Mode)
After=network.target readsb.service gpsd.service
Wants=gpsd.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/adsb-sensor

# -----------------------------------------------------------
# HARDWARE CONFIGURATION
# -----------------------------------------------------------
Environment="SENSOR_ID={{ sensor_id | default(inventory_hostname) }}"
Environment="SENSOR_GNSS_PATH={{ gnss_device_path | default('/dev/ttyUSB0') }}"
Environment="PYTHONUNBUFFERED=1"

# -----------------------------------------------------------
# EXECUTION COMMAND (FIXED)
# -----------------------------------------------------------
# Now dynamically adds flags based on Ansible inventory variables
ExecStart=/usr/bin/python3 /opt/adsb-sensor/smart_adsb_logger.py \
    --host localhost \
    --port 30005 \
    --log-dir /var/lib/adsb_storage/csv_data \
    {% if log_aircraft_csv | default(false) %}--log-aircraft{% endif %}

# RESTART POLICY
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
