# infra/ansible/roles/gnss-receiver/templates/gnss-receiver.service.j2
[Unit]
Description=GNSS Receiver (Socat TCP Hub)
After=network.target
# This ensures our init script ran first
Requires=gnss-receiver.service

[Service]
Type=simple
User=root
Restart=always
RestartSec=5

# CRITICAL: We add ',b230400' to ensure socat keeps the high speed
# and doesn't revert the port to 9600 defaults.
ExecStart=/usr/bin/socat -d -d -v TCP-LISTEN:2000,fork,reuseaddr FILE:/dev/ttyACM0,raw,echo=0,b230400

[Install]
WantedBy=multi-user.target
