# infra/ansible/roles/gnss-receiver/templates/gnss-ntrip-client.service.j2
# ==============================================================================
# UNIT: NTRIP Client (RTK Injector)
# FILE: infra/ansible/roles/gnss-receiver/templates/gnss-ntrip-client.service.j2
# VERSION: 0.4.4
# ==============================================================================

[Unit]
Description=NTRIP Client (RTK Injector)
After=gnss-socat-hub.service network-online.target
Wants=gnss-socat-hub.service

[Service]
Type=simple
User=root
Restart=always
RestartSec=15

# ACADEMIC LOGIC:
# 1. Input: NTRIP Caster (Remote)
# 2. Output: Local Data Hub (tcpcli://localhost:{{ socat_hub_port }})
#    - This pushes RTCM corrections INTO the socat hub, which writes them to HW.

ExecStart=/usr/local/bin/str2str \
  -in ntrip://{{ ntrip_user }}:{{ ntrip_password }}@{{ ntrip_host }}:{{ ntrip_port }}/{{ ntrip_mountpoint }} \
  -out tcpcli://127.0.0.1:{{ socat_hub_port }} \
  -msg "1002,1006,1013,1019,1033,1077,1087,1097,1107,1127,1230" \
  -n 10

[Install]
WantedBy=multi-user.target
