# infra/ansible/roles/gnss-receiver/templates/gnss-receiver.service.j2
[Unit]
Description=GNSS Hardware Hub (str2str)
# Description: Reads raw Serial Data and exposes it on TCP 2000.
# This is the "Root" of the RTK data pipeline.
After=network.target
Before=gnss-splitter.service

[Service]
Type=simple
Restart=always
RestartSec=5
# Run initialization script before starting the stream
ExecStartPre=/usr/local/bin/ublox_init.sh
# Start str2str:
# -in serial://ttyACM0:230400  -> Read from USB at High Speed
# -out tcpsvr://:2000          -> Serve raw data on Local Port 2000
ExecStart=/usr/local/bin/str2str -in serial://ttyACM0:230400#ubx -out tcpsvr://:2000
ExecStop=/usr/bin/killall str2str

[Install]
WantedBy=multi-user.target
