# infra/ansible/roles/gnss-receiver/templates/gnss-configurator.service.j2
[Unit]
Description=GNSS Hardware Configurator (Silence NMEA)
# Runs once at boot to force the chip into Binary Mode
After=gnss-receiver.service
Requires=gnss-receiver.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/configure_gnss.sh
RemainAfterExit=yes
User=root
Group=root

[Install]
WantedBy=multi-user.target
